# services/crypt_boot_routine.service
[Unit]
Description=Run decryption at login

[Service]
Type=oneshot
ExecStart=
Restart=no

[Install]
WantedBy=default.target